# codebaux/guest/overlay/etc/init.d/codebaux-serial
#!/sbin/openrc-run
# ─────────────────────────────────────────────────────────────
# codebaux-serial
# OpenRC service that starts the serial-listener daemon.
# This is the entry point for all host↔guest communication.
# ─────────────────────────────────────────────────────────────

name="codebaux-serial"
description="Codebaux serial protocol listener"
command="/usr/local/bin/serial-listener"
command_background="yes"
pidfile="/run/${name}.pid"
output_log="/var/log/codebaux-serial.log"
error_log="/var/log/codebaux-serial.log"

depend() {
    need localmount
    after bootmisc
}

start_pre() {
    # Ensure workspace directories exist
    mkdir -p /workspace /workspace/out
    # Ensure log directory
    mkdir -p /var/log
    touch /var/log/project.log
}
